# server/com.sendtoopenclaw.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sendtoclwdbot</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/opt/send-to-openclaw/server/server.js</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/opt/send-to-openclaw/server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PORT</key>
    <string>3847</string>
    <key>OPENCLAW_WAKE_MODE</key>
    <string>now</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>
